# skeleton-0.0.0/skeleton/etc/inittab
null::sysinit:/bin/mount -t proc proc /proc
null::sysinit:/bin/mount -o remount,rw /
null::sysinit:/sbin/busybox --install /sbin
null::sysinit:/sbin/busybox --install /bin
null::sysinit:/bin/mkdir -p /dev/pts
null::sysinit:/bin/mkdir -p /dev/shm
null::sysinit:/bin/mount -a
null::sysinit:/bin/mount -t devpts devpts /dev/pts
null::sysinit:/bin/hostname -F /etc/hostname
tty1::respawn:/sbin/getty -L tty1 115200 vt100
null::shutdown:/bin/umount -a -r
null::shutdown:/bin/swapoff -a
null::sysinit:/sbin/chmod 4755 /sbin/busybox
null::sysinit:/etc/rc.d/udhcpc start
null::sysinit:/etc/rc.d/dropbear start
null::sysinit:/etc/rc.d/syslogd start
null::sysinit:/etc/rc.d/klogd start
